# templates/indeximas.service
[Unit]
Description=Indexima Service
After=network.target
StartLimitInterval=60
StartLimitBurst=5

[Service]
EnvironmentFile=/etc/environment
Restart=always
WorkingDirectory={{ galactica_path }}
LimitNOFILE=65536
User={{ service_user | d('root') }}
ExecStart=/bin/bash {{ galactica_path }}/start-node.sh

[Install]
WantedBy=multi-user.target